# supervisord.conf
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[supervisord]
pidfile=supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[program:gunicorn]
command=gunicorn hydra.wsgi:application -c production.py
autostart=false